# 4_Backup/backup-0.1/backup.service
Unit]
Description=Backup service

[Service]
Type=oneshot
ExecStart=/usr/bin/mybackup.sh
User=#USER#

[Install]
WantedBy=multi-user.target
