# 4_Backup/backup-0.1/backup.timer
[Unit]
Description=my backup

[Timer]
OnCalendar=*-*-* 01:00:00
Unit=backup.service

[Install]
WantedBy=multi-user.target

